# rasp5/tools/GuiPnCrema.service
[Unit]
  Description=PallanuotoCrema GUI (user session)
  After=graphical-session.target network-online.target
  # IMPORTANTE: Aspetta che la sessione grafica sia pronta
  Wants=graphical-session.target

  [Service]
  Type=simple
  WorkingDirectory=/home/davide/segnapunti_pallanuoto/
  ExecStart=/home/davide/python_venv/bin/python /home/davide/segnapunti_pallanuoto/rasp5/gui/giuria_gui.py
  Restart=on-failure
  RestartSec=5

  # Variabili d'ambiente per Qt/X11
  Environment=DISPLAY=:0
  Environment=XDG_RUNTIME_DIR=/run/user/1000
  Environment=XAUTHORITY=/home/davide/.Xauthority
  Environment=QT_QPA_PLATFORM=xcb

  [Install]
  WantedBy=graphical-session.target
